# qr3_92.conf
[program:qr3_92]
command=python3 /home/curator3/qr3_92/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_92/test.err.log
stdout_logfile=/home/curator3/qr3_92/test.err.log
